# etc/systemd/system/steamos-exit-splash.service
[Unit]
# =============================================================================
# SteamMachine-DIY - Exit Splash Screen Service
# =============================================================================
Description=SteamOS-DIY Exit Splash
DefaultDependencies=no
# We must run before the targets that actually stop services
Before=shutdown.target reboot.target halt.target
# Conflicts with umount.target ensures we stay alive until the very end
Conflicts=umount.target

[Service]
Type=oneshot
# Logic: We detect the target directly in a single shell execution.
# This avoids running the splash twice and is much faster.
ExecStart=/bin/sh -c '\
    if systemctl list-jobs | grep -q "reboot.target"; then \
        /usr/local/bin/steamos-helpers/steamos-splash reboot; \
    else \
        /usr/local/bin/steamos-helpers/steamos-splash shutdown; \
    fi'

StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty1
# Keep the service "active" so it's not cleaned up too early
RemainAfterExit=yes

[Install]
WantedBy=shutdown.target reboot.target
